# etc/systemd/transport.service
[Unit]
Description=Мониторинг транспорта накладных
 
[Service]
WorkingDirectory=/home/anton/script/AveWd
ExecStart=/bin/sh transport.sh


