# etc/systemd/transport.timer
[Unit]
Description=Мониторинг транспорта накладных

[Timer]
OnBootSec=1min
OnUnitActiveSec=3hours

[Install]
WantedBy=timers.target